# supervisor/website.conf
[program:uwsgi-app]
command = /usr/bin/uwsgi-core --ini /home/uwsgi/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx